# test/dfuzzer-test-server.service.in
[Unit]
Description=dfuzzer test server

[Service]
Type=dbus
BusName=org.freedesktop.dfuzzerServer
ExecStart=@BINDIR_RAW@/dfuzzer-test-server
