# src/dfuzzer-test-server.service.in
[Unit]
Description=dfuzzer test server

[Service]
Type=dbus
BusName=org.freedesktop.dfuzzerServer
ExecStart=@PREFIXDIR_RAW@/bin/dfuzzer-test-server
